# etc/uci-defaults/set_zsh_as_default_shell
which zsh && sed -i -- 's:/bin/ash:'`which zsh`':g' /etc/passwd
